# etc/services.d/websockify/run
#!/usr/bin/with-contenv bash
set -e
exec websockify --web=/usr/share/novnc/ 6080 localhost:5900